# week04_automation_jobs_scripting/day25_alerts_logs/monitoring.cron
# Run every day at 7:00 AM
0 7 * * * /path/to/run.sh >> /path/to/output/monitor_log.txt 2>&1


crontab cron/monitoring.cron
